# etc/inittab
::sysinit:/etc/init.d/rc
tty0::askfirst:/bin/cttyhack /bin/sh
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
#tty0::respawn:-/bin/sh
ttyS0::respawn:-/bin/sh